# klogd.service.in
[Unit]
Description=Kernel Logging Service
After=syslogd.service

[Service]
Type=simple
ExecStart=@SBINDIR@/klogd -n
StandardOutput=null
Restart=on-failure

[Install]
WantedBy=multi-user.target
